# app-crypt/uacme/files/ualpn.initd
#!/sbin/openrc-run
command=/usr/bin/ualpn
pidfile="/run/${RC_SVCNAME}-svc.pid"
command_args="${ualpn_args}"
command_args_background="-d --pidfile ${pidfile}"

name="ualpn"
description="tls-alpn-01 challenge responder"

depend() {
        need net
        use logger
        before nginx
}
