# app-crypt/uacme/files/ualpn.initd
#!/sbin/openrc-run
command=/usr/bin/ualpn
pidfile="/run/${RC_SVCNAME}.pid"
command_args="--pidfile ${pidfile} ${ualpn_args}"
command_args_background="-d"

name="ualpn"
description="tls-alpn-01 challenge responder"

depend() {
        need net
        use logger
        before nginx
}
