# app-metrics/gigahost-agent/files/gigahost-agent.cron
* * * * * ghagent bash /opt/gigahost/agent.sh > /opt/gigahost/cron.log 2>&1
